# ops/crontab
0 * * * * /home/zdv/vzdv/ops/backup.sh >> /home/zdv/vzdv/backup.log 2>&1
